# seleniumAutomation.service
[Unit]
Description=Start MailHog and Selenium Automation

# Start ONLY after these are fully started
After=network-online.target docker.service jsServer.service apache2.service
Requires=jsServer.service apache2.service
Wants=network-online.target docker.service

[Service]
Type=simple
WorkingDirectory=/home/admin
User=admin
Group=admin  # Assuming 'admin' group exists; adjust if needed

# Wait for network to be online (more robust check)
ExecStartPre=/bin/bash -c 'until ping -c1 -W1 8.8.8.8 > /dev/null 2>&1; do sleep 2; done'

# Wait for Docker daemon to be ready
ExecStartPre=/bin/bash -c 'until docker info > /dev/null 2>&1; do sleep 2; done'

# Wait for jsServer to be listening (assuming it runs on a port; replace 8080 with actual port if known)
ExecStartPre=/bin/bash -c 'until curl -sf http://127.0.0.1:8080 > /dev/null; do sleep 2; done'

# Wait for Apache to be listening (your original check)
ExecStartPre=/bin/bash -c 'until curl -sf http://127.0.0.1/ > /dev/null; do sleep 2; done'

ExecStart=/home/admin/csrfScenario/automation.sh
KillMode=control-group
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target