# services/restartServer.service
[Unit]
Description=Restart backend on enp0s3 IP change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart jsServer.service

