# services/seleniumAutomation.service
[Unit]
Description=Start Selenium Automation

# Start ONLY after these are fully started
After=network-online.target jsServer.service apache2.service
Requires=jsServer.service apache2.service
Wants=network-online.target

[Service]
Type=simple
# No blind sleep: autoLogin.py waits for the app (frontend + API) to actually
# respond via wait_for_app_ready() before attempting login.
ExecStart=/opt/selenium_venv/bin/python3 /usr/local/bin/autoLogin.py
KillMode=control-group
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
